# deploy/sovereign-explorer.service
# systemd unit for the Sovereign block explorer on a small VPS.
#
# Install:
#   sudo useradd --system --create-home --home-dir /opt/sov-explorer explorer  # or reuse
#   sudo git clone https://github.com/cloudzombie/sov-explorer /opt/sov-explorer/app
#   sudo install -m 0600 deploy/sov-explorer.env.example /etc/sov-explorer.env
#   # Provision the TLS relay names and edit /etc/sov-explorer.env, then:
#   sudo cp /opt/sov-explorer/app/deploy/sovereign-explorer.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now sovereign-explorer
#   journalctl -u sovereign-explorer -f
#
# The archive is durable but derived entirely from identity-pinned relays, so it can
# be rebuilt. Bind to loopback and publish through the supplied TLS reverse proxy.

[Unit]
Description=Sovereign block explorer
Documentation=https://github.com/cloudzombie/sov-explorer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=explorer
Group=explorer
WorkingDirectory=/opt/sov-explorer/app
ExecStartPre=/usr/bin/node -e "const [a,b]=process.versions.node.split('.').map(Number);if(a<24||(a===24&&b<15)){console.error('Node 24.15+ required');process.exit(1)}"
ExecStart=/usr/bin/node src/server.js

Environment=PORT=8730
Environment=HOST=127.0.0.1
Environment=ARCHIVE_DIR=/var/lib/sovereign-explorer
Environment=ARCHIVE_BACKFILL_BATCH=16
Environment=REQUIRE_TLS_RELAYS=1
Environment=TRUST_PROXY=1
Environment=HTTP_REQUESTS_PER_MINUTE=600
Environment=GRAPHQL_REQUESTS_PER_MINUTE=60
Environment=GLOBAL_HTTP_REQUESTS_PER_MINUTE=30000
Environment=GLOBAL_GRAPHQL_REQUESTS_PER_MINUTE=3000
Environment=HTTP_MAX_CONNECTIONS=2000
Environment=WS_MAX_CLIENTS=1000
Environment=WS_MAX_CLIENTS_PER_NETWORK=750
# Required: the explorer fails closed until real HTTPS relay endpoints are supplied.
EnvironmentFile=/etc/sov-explorer.env

Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sovereign-explorer

# Service sandboxing. Node/V8 needs executable JIT memory, so MemoryDenyWriteExecute
# is intentionally not enabled.
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
UMask=0077
StateDirectory=sovereign-explorer
StateDirectoryMode=0700
LimitNOFILE=65536
TasksMax=128
MemoryHigh=640M
MemoryMax=768M

[Install]
WantedBy=multi-user.target
